# com_goldenthinker_trade_install/ec2/etc/systemd/system/goldenthinker_services/gt_telegram_incoming_signal.service
[Unit]
Description=Telegram Incoming Signal Service
After=multi-user.target

[Service]
Type=simple
WorkingDirectory=/crypto_trading
ExecStart=/usr/local/bin/python /crypto_trading/gt_telegram_incoming_signal.py

[Install]
WantedBy=multi-user.target